# inf/crontab/cron
*/5 * * * * . /home/ubuntu/app/env/.env; java -jar /home/ubuntu/app/prod/cron/wtc-backend-cli.jar POST_LOADER
