# inf/crontab/cron
*/5 * * * * java -jar /home/ubuntu/app/prod/cron/wtc-backend-cli.jar POST_LOADER
