# inf/crontab/cron
*/5 * * * * /root/app/cron/wtc-post-loader
