# inf/crontab/cron
*/5 * * * * java -jar /root/app/prod/cron/wtc-backend-cli.jar POST_LOADER
